# supervisord.d/c9.ini
[program:c9]
command=/root/.c9/node/bin/node /usr/src/c9sdk/server.js -p 444 -a admin:admin --packed
autorestart = true