# supervisord.d/c9.ini
[program:c9]
directory=/usr/local/src/c9sdk
command=/root/.c9/node/bin/node /usr/local/src/c9sdk/server.js -p 444 -a admin:THE_PASSWORD -w /root/workspace --packed
autorestart = true