# supervisord.d/caddy.ini
[program:caddy]
command=/usr/local/bin/caddy -conf /etc/caddy/Caddyfile
autorestart = true