# supervisord.d/websockify.ini
[program:websockify]
directory=/usr/local/src/websockify
command=/usr/local/src/websockify/run 82 127.0.0.1:5901
autorestart = true